# deploy/gateway.service
[Unit]
Description=Gateway Service
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/gateway
ExecStart=/home/ubuntu/gateway/gateway
Restart=always
RestartSec=5
Environment=RUST_LOG=info

# 从 .env 文件加载环境变量
EnvironmentFile=/home/ubuntu/gateway/.env

[Install]
WantedBy=multi-user.target
